# dist/systemd/rq-dashboard.service
[Unit]
Description=RQ Dashboard
Requires=redis-server.service gunicorn.service
After=redis-server.service gunicorn.service
BindsTo=gunicorn.service

[Service]
Type=simple
User=antbs
Group=antbs
ExecStart=/usr/bin/rq-dashboard -p 9181 -b 127.0.0.1
TimeoutStopSec=10
WorkingDirectory=/PATH/TO/antbs/antbs
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
